# backend/supervisord.conf
; Runs the Exposure Triage Assistant backend under supervisor — no root
; required. Port 8001, deliberately different from gbs-fastapi's 8000, since
; both run on the same host at once. All paths below use %(here)s, which
; supervisor expands to the directory containing this file, so this config
; works regardless of where the repo is checked out.
;
; Setup:
;   pip install -r requirements.txt
;   .venv/bin/supervisord -c supervisord.conf
;
; Day to day:
;   .venv/bin/supervisorctl -c supervisord.conf status
;   .venv/bin/supervisorctl -c supervisord.conf restart vulnerabilityscanner-backend
;   .venv/bin/supervisorctl -c supervisord.conf tail -f vulnerabilityscanner-backend
;
; Survive a server reboot (no systemd needed): add a crontab @reboot entry —
; see the README's Deployment section.

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisord]
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=5
pidfile=%(here)s/supervisord.pid
childlogdir=%(here)s/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[program:vulnerabilityscanner-backend]
command=%(here)s/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001
directory=%(here)s
environment=PATH="%(here)s/.venv/bin:%(ENV_PATH)s"
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stdout_logfile=%(here)s/logs/uvicorn.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/logs/uvicorn.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=5
